# oi-v1/deploy/oi-server.service
[Unit]
Description=oi desk-device HTTP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dev
WorkingDirectory=/opt/oi
Environment=OI_PORT=8842
Environment=OI_STATE_DIR=/var/lib/oi
ExecStart=/usr/bin/python3 /opt/oi/server/oi_server.py --host 0.0.0.0 --port ${OI_PORT} --state-dir ${OI_STATE_DIR}
Restart=on-failure
RestartSec=2
StateDirectory=oi

# Basic hardening. State is written only under /var/lib/oi.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/oi

[Install]
WantedBy=multi-user.target
